# sentinel/crontab
* * * * * cd /home/vivo/sentinel && SENTINEL_CONFIG=/home/vivo/sentinel/sentinel.conf python /opt/wallet/vivo/sentinel/bin/sentinel.py
